# configs/.crontab
0 0 */1 * * /Users/clarkbenham/screenfilm/start.sh
